# ansible-playbooks/backup/da-x-posts-backup/files/watch_sftp_dir.service
# Trigger event when upload by SFTP
#

[Unit]
Description=trigger event for SFTP

[Service]
ExecStart=/usr/local/sbin/watch_sftp_dir.sh
