# ansible-playbooks/databases/mysql-base/files/daily-backup.timer
# daily database backup timer
#

[Unit]
Description=Daily Database Backup Timer

[Timer]
Unit=daily-backup.service
OnCalendar=*-*-* 22:15:00
Persistent=true

[Install]
WantedBy=default.target
