# ansible-playbooks/databases/mysql-base/files/daily-backup.service
# daily database backup service call
#

[Unit]
Description=Daily database backup service call

[Service]
ExecStart=/usr/local/sbin/daily-backup.sh
